# docker/supervisor.d/services.ini
[program:database]
command=/usr/bin/postgres -D /db
redirect_stderr=true
stdout_logfile=/db/db.log
user=postgres
autostart=false

[program:serve]
command=/usr/bin/python3 -u /app/serve.py
directory=/app
redirect_stderr=true
stdout_logfile=/run/app.log
user=app
autostart=false